# opt/upgrade/2.0_to_2.5/docs/crontab
Insert the following into your crontab file

*/5 * * * * /opt/nessquik-client/scripts/scan-runner.php >/dev/null 2>&1
